# systemd/phoniebox-hotspot-fallback.service
[Unit]
Description=Phoniebox Hotspot Fallback Cycle
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/phoniebox-panel
ExecStart=/opt/phoniebox-panel/.venv/bin/python /opt/phoniebox-panel/scripts/fallback_hotspot.py
