# systemd/phoniebox-gpio-poll.service
[Unit]
Description=Phoniebox GPIO Polling
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/phoniebox-panel
ExecStart=/opt/phoniebox-panel/.venv/bin/python /opt/phoniebox-panel/scripts/gpio_poll.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
