# systemd/phoniebox-network-bootstrap.service
[Unit]
Description=Bootstrap Phoniebox WLAN and hotspot profiles
After=NetworkManager.service network-online.target
Wants=NetworkManager.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/phoniebox-panel
ExecStart=/opt/phoniebox-panel/.venv/bin/python /opt/phoniebox-panel/scripts/bootstrap_network.py --seed-only

[Install]
WantedBy=multi-user.target
