# systemd/phoniebox-runtime-tick.service
[Unit]
Description=Phoniebox Runtime Tick
After=network.target phoniebox-panel.service
Wants=phoniebox-panel.service

[Service]
Type=simple
WorkingDirectory=/opt/phoniebox-panel
ExecStart=/opt/phoniebox-panel/.venv/bin/python /opt/phoniebox-panel/scripts/runtime_tick.py --daemon --interval 1.0
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
