# projects/Phoniebox/panel/systemd/phoniebox-runtime-tick.service
[Unit]
Description=Phoniebox Runtime Tick
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/phoniebox-panel
ExecStart=/usr/bin/python3 /opt/phoniebox-panel/scripts/runtime_tick.py
